# ops/native/plists/com.robinco.auto-trader.mcp-tradingcodex-execution.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.robinco.auto-trader.mcp-tradingcodex-execution</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/mgh3326/services/auto_trader/scripts/run-mcp-profile.sh</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/mgh3326/services/auto_trader/current</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/mgh3326</string>
    <key>PATH</key>
    <string>/Users/mgh3326/.local/bin:/Users/mgh3326/.hermes/node/bin:/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>AUTO_TRADER_BASE</key>
    <string>/Users/mgh3326/services/auto_trader</string>
    <key>AUTO_TRADER_ENV_FILE</key>
    <string>/Users/mgh3326/services/auto_trader/shared/.env.prod.native</string>
    <key>AUTO_TRADER_MCP_PROFILE</key>
    <string>tradingcodex_execution</string>
    <key>AUTO_TRADER_MCP_PORT</key>
    <string>8770</string>
    <key>AUTO_TRADER_MCP_AUTH_TOKEN_ENV</key>
    <string>MCP_TRADINGCODEX_EXECUTION_AUTH_TOKEN</string>
    <key>ORDER_APPROVAL_HASH_MODE</key>
    <string>required</string>
    <key>TOSS_APPROVAL_HASH_MODE</key>
    <string>required</string>
  </dict>
  <key>RunAtLoad</key>
  <false/>
  <key>KeepAlive</key>
  <true/>
  <key>ThrottleInterval</key>
  <integer>10</integer>
  <key>StandardOutPath</key>
  <string>/Users/mgh3326/services/auto_trader/logs/com.robinco.auto-trader.mcp-tradingcodex-execution.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/mgh3326/services/auto_trader/logs/com.robinco.auto-trader.mcp-tradingcodex-execution.err.log</string>
</dict>
</plist>
